# config/supervisor/prod.ini
[program:song_manager]
process_name=song_manager_%(process_num)s
environment=NODE_ENV=production,PORT=900%(process_num)s,CONFIG_FILE="/home/song-manager/config/song-manager.json"
command=node src/server.js
autostart=true
autorestar=true
numprocs=2
