# config/supervisor/prod.ini
[program:song_manager]
process_name=song_manager_%(process_num)s
environment=NODE_ENV=production,PORT=900%(process_num)s,CONFIG_FILE="/home/song-manager/config/song-manager.json"
directory=/opt/apps/song-manager
command=node src/server.js
user=song-manager
autostart=true
autorestar=true
numprocs=2
stdout_logfile=/var/log/song-manager/app-worker-%(process_num)s.log
stdout_logfile_maxbytes=50MB
stdout_backups=10
stderr_logfile=/var/log/song-manager/error-worker-%(process_num)s.log
stderr_logfile_maxbytes=50MB
stderr_backups=10
